# roles/host-setup/templates/docker.service.d/override.conf.j2
[Unit]
After=network-online.target ufw.service containerd.service

[Service]
ExecStartPre=/bin/sleep 2
